# deploy/systemd/hermes-mail-agent.service
[Unit]
Description=OpenClaw Multimodal Mail Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/openclaw-multimodal-toolkit
EnvironmentFile=/opt/openclaw-multimodal-toolkit/hermes-mail-agent.env
ExecStart=/usr/bin/python3 /opt/openclaw-multimodal-toolkit/mail-agent.py
StandardOutput=journal
StandardError=journal
MemoryMax=512M
CPUQuota=50%
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
